# templates/com.user.ralph-crew.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.ralph-crew</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/dotfiles/scripts/ralph-crew</string>
        <string>dispatch</string>
        <string>--config</string>
        <string>__PROJECT__/.claude/crew.json</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>
    <key>StartInterval</key>
    <integer>__INTERVAL__</integer>
    <key>StandardOutPath</key>
    <string>/tmp/ralph-crew/logs/launchd.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ralph-crew/logs/launchd.err</string>
</dict>
</plist>
