# templates/com.user.claude-status-local-watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.claude-status-local-watch</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/dotfiles/scripts/claude-status-local-watch.sh</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>/tmp/claude_status</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/tmp/claude-status-local-watch.err</string>
    <key>StandardOutPath</key>
    <string>/tmp/claude-status-local-watch.log</string>
</dict>
</plist>
